# service/horseracing-server.service
[Unit]
Description=Horser Racing REST API Server
Before=multi-user.target

[Service]
WorkingDirectory=/work/backend
ExecStart=/usr/bin/python3 ./server.py --start
ExecStop=/usr/bin/python3 ./server.py --stop
Type=simple
Restart=always

[Install]
WantedBy=default.target
RequiredBy=network.target