# service/horseracing-socket.service
[Unit]
Description=Horser Racing Socket Server
Before=multi-user.target

[Service]
WorkingDirectory=/work/backend
ExecStart=/usr/bin/python3 ./socket-server.py --start
ExecStop=/usr/bin/python3 ./socket-server.py --stop
Type=simple
Restart=always

[Install]
WantedBy=default.target
RequiredBy=network.target